# config/robert-ia.service
[Unit]
Description=Robert-IA Backend
After=network-online.target local-fs.target
Wants=network-online.target

[Service]
Type=simple
User=robert-ia
Group=robert-ia
WorkingDirectory=/opt/robert-ia
ExecStart=/opt/robert-ia/scripts/start-backend.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=robert-ia
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/home/robert-ia"
CPUQuota=80%
MemoryLimit=3G
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
